# sv/pipewire/run
#!/bin/sh

exec 2>&1

pkill --signal KILL pulseaudio
pkill --signal KILL pipewire
pkill --signal KILL wireplumber

exec pipewire -v
