# sv/pipewire-pulse/run
#!/bin/sh

exec 2>&1
SVDIR="$HOME/sv" "$HOME/bin/svwaitup" pipewire || (sleep 60; exit 1)
exec pipewire-pulse
