# sv/syncthing/run
#!/bin/sh

: ${HOME:="$(eval print "~$(id -un)")"}
: ${XDG_CONFIG_HOME:="$HOME"/.config}

exec 2>&1
exec syncthing --no-browser
