# profile.d/colors
#!/bin/sh
# vim: ft=sh ts=2 sw=2 et:

# GNU ls
# TODO: Replace with third party LS_COLORS tool.
for b in \
  /usr/local/opt/coreutils/libexec/gnubin/dircolors    \
  /opt/homebrew/opt/coreutils/libexec/gnubin/dircolors \
  "$(command -v dircolors 2>/dev/null)"; do
  [ -x "$b" ] && eval "$("$b" -b "$CONF/dir_colors")"
done
unset b

# General
# https://bixense.com/clicolors/
export CLICOLOR=1

# Apps
export HIGHLIGHT_STYLE="nord"
