# sv/pmd/run
#!/bin/sh

: ${HOME:="$(eval print "~$(id -un)")"}

exec 2>&1
exec "$HOME"/bin/pmd
