# sv/picom/run
#!/bin/sh

: ${HOME:="$(eval print "~$(id -un)")"}
: ${XDG_CONFIG_HOME:="$HOME"/.config}

exec 2>&1
exec picom --config "$XDG_CONFIG_HOME/picom/picom.conf" --experimental-backends
