# sv/1password/run
#!/bin/sh

exec 2>&1

1password --disable-smooth-scrolling --silent
exec pause
