# sv/transmission/run
#!/bin/sh

exec 2>&1
exec transmission-daemon \
  --foreground \
  --log-level=info --logfile /dev/stdout \
  --watch-dir "$HOME/downloads/torrents" \
  --incomplete-dir "$HOME/downloads/pending"
