# sv/weechat/run
#!/bin/sh

: ${HOME:="$(eval print "~$(id -un)")"}
: ${XDG_DATA_HOME:="$HOME"/.local}

exec 2>&1
exec chpst -b weechatd dtach -N "$XDG_DATA_HOME/dtach/weechat" weechat
