# sv/i3/run
#!/bin/sh

exec 2>&1
cd
exec i3 -c "$XDG_CONFIG_HOME/i3/main.conf"
